# rootfs/etc/systemd/system/firewalld.service.d/override.conf
[Service]
# ExecStartPre=/warewulf/sbin/execstartpre.%N.sh
# ExecStartPost=/warewulf/sbin/execstartpost.%N.sh
TimeoutStartSec=900

